# code/supervisord.conf
# supervisord.conf
[program:goji]
command=/usr/local/bin/einhorn -c goji -b 127.0.0.1:1234 -m manual example -bind einhorn@0
